# com.simone.kb-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.simone.kb-manager</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>/Users/simone/dev_ai/kb_manager/launcher_server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/simone/dev_ai/kb_manager</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/simone/dev_ai/kb_manager/server.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/simone/dev_ai/kb_manager/server.log</string>
</dict>
</plist>
